# config/systemd/system/kitchenradio-fip-stop.timer
[Unit]
Description=Stop timer for radio

[Timer]
OnActiveSec=4h
Unit=kitchenradio-fip.service

[Install]
WantedBy=timers.target
