# config/systemd/system/kitchenradio-fip-stop.service
[Unit]
Description=Stop radio

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop kitchenradio-fip.service

[Install]
WantedBy=default.target
